# slurm/s6-rc.d/miner-api/run
#!/bin/sh

set -e
/usr/sbin/miner-api -D
